# launchagent/com.sloan.agentos.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.sloan.agentos.refresh</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>/Users/jes/AI_Agent_OS/bin/agentos_refresh.sh</string>
	</array>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>900</integer>
	<key>WorkingDirectory</key>
	<string>/Users/jes/AI_Agent_OS</string>
	<key>StandardOutPath</key>
	<string>/Users/jes/AI_Agent_OS/reports/launchagent.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/jes/AI_Agent_OS/reports/launchagent.error.log</string>
</dict>
</plist>
